# scripts/services/frontend.service
[Unit]
Description=Frontend Service
After=network.target

[Service]
Type=simple
User=STRAMATEL
WorkingDirectory=/home/STRAMATEL/builds/builds/
ExecStart=/bin/bash /home/STRAMATEL/builds/builds/startFrontend.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

